# scripts/fa-rk3328-misc
#!/bin/sh /etc/rc.common

START=96
 
start() {
    [ -f /proc/irq/29/smp_affinity ] && echo 8 > /proc/irq/29/smp_affinity
}
